# backend/celery.service
[Unit]
Description=Celery Worker Service for Librarity
After=network.target redis.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/Librarity/backend
Environment="PATH=/root/Librarity/backend/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/root/Librarity/backend/venv"
Environment="PYTHONPATH=/root/Librarity/backend"
ExecStart=/root/Librarity/backend/venv/bin/python -m celery -A workers.celery_app worker --loglevel=info --pool=solo --concurrency=1

ExecStop=/bin/kill -s TERM $MAINPID
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=always
RestartSec=10

# Security
PrivateTmp=true

[Install]
WantedBy=multi-user.target
